# meta-shadowhunt/recipes-core/studio-system/files/display.cron
# reload the webpage every 5 minutes
DISPLAY=:0

*/5 * * * * root /usr/bin/dbus-send --print-reply --session --dest=de.shadowhunt.display /browser de.shadowhunt.display.BrowserService.reload
